# dev-stories/automation/dev-story-preview-comment-current.service
[Unit]
Description=dev-story: every open PR has a current preview-link comment
After=network-online.target

[Service]
Type=oneshot
User=testuser
Environment=HOME=/home/testuser
ExecStart=/home/testuser/bin/dev-story-preview-comment-current.sh
