# dev-stories/automation/dev-story-demo-envs-reachable.service
[Unit]
Description=dev-story: demo environments are reachable and live
After=network-online.target

[Service]
Type=oneshot
User=testuser
Environment=HOME=/home/testuser
ExecStart=/home/testuser/bin/dev-story-demo-envs-reachable.sh
