# dev-stories/automation/dev-story-deploy-freshness.service
[Unit]
Description=dev-story: production reflects the latest merge (deploy freshness)
After=network-online.target

[Service]
Type=oneshot
User=testuser
Environment=HOME=/home/testuser
ExecStart=/home/testuser/bin/dev-story-deploy-freshness.sh
